# configs/automatic-backup/orangehrm-deployment-backup.service
[Unit]
Description=OrangeHRM Deployment Backup

[Service]
ExecStart=/usr/bin/orangehrm-deployment-backup.sh
